# deploy/com.local.hob.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Hob as a user LaunchAgent. Copy to ~/Library/LaunchAgents/, edit the /Users/you
     paths and HOB_ALLOWED_TELEGRAM_USER_ID. Store the token first with:
       uv run python app.py token set
     Then:
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.local.hob.plist
       launchctl kickstart -k gui/$(id -u)/com.local.hob       # (re)start
-->
<plist version="1.0">
<dict>
  <key>Label</key>            <string>com.local.hob</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/you/.local/bin/uv</string>
    <string>run</string>
    <string>--directory</string>
    <string>/Users/you/Git/Hob</string>
    <string>python</string>
    <string>app.py</string>
  </array>
  <key>WorkingDirectory</key>  <string>/Users/you/Git/Hob</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOB_ALLOWED_TELEGRAM_USER_ID</key> <string>123456789</string>
    <key>HOB_MODEL</key>          <string>qwen2.5:14b-instruct</string>
    <key>HOB_WAKE_TIME</key>      <string>07:00</string>
    <key>HOB_DB_PATH</key>        <string>/Users/you/Library/Application Support/Hob/hob.db</string>
    <key>HOB_KEEP_ALIVE</key>     <string>-1</string>
    <key>HOB_CALENDAR_ENABLED</key><string>true</string>
    <key>HOB_WORK_HOURS</key>      <string>09:00-17:30</string>
    <key>HOB_WORK_DAYS</key>       <string>mon,tue,wed,thu,fri</string>
    <key>HOB_BREAKS</key>          <string>12:00-13:00</string>
    <key>HOB_DEFAULT_DURATION</key><string>30</string>
    <key>HOB_TRANSITION_BUFFER</key><string>0</string>
    <key>PATH</key>              <string>/Users/you/.local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>        <true/>
  <key>KeepAlive</key>        <true/>
  <key>ProcessType</key>      <string>Standard</string>
  <key>StandardOutPath</key>  <string>/Users/you/Library/Application Support/Hob/hob.log</string>
  <key>StandardErrorPath</key><string>/Users/you/Library/Application Support/Hob/hob.log</string>
</dict>
</plist>
